# dotfiles/cachyos-cinnamon/.config/autostart/numlock.desktop
[Desktop Entry]
Name=Enable NumLock
Comment=Enable NumLock on startup
Exec=numlockx on
Type=Application
X-GNOME-Autostart-enabled=true

